# hunt-ui/daemon/com.ayaan.paperhunt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ayaan.paperhunt</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/ayaansmacmini2/.paperhunt/run-hunt.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ayaansmacmini2/.paperhunt/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ayaansmacmini2/.paperhunt/logs/launchd.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ayaansmacmini2/.paperhunt</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ayaansmacmini2/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ayaansmacmini2</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
